# Twitter/InternetIsRealLife/internet.cron
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
00 2,6,10,14,18,22 * * * bash /home/orange/Projects/Scripts/Twitter/InternetIsRealLife/PhotoTweet.sh > /home/orange/Projects/Scripts/cronoutput

